# mopidy-local-scan.service
[Unit]
Description=Scan local mopidy database

[Service]
ExecStart=/usr/sbin/mopidyctl local scan

[Install]
WantedBy=basic.target
